# templates/systemd/fleet-heartbeat.service
[Unit]
Description=Ollie Fleet heartbeat (push stack health to the fleet dashboard)
# No-ops cleanly until ~/.config/ollie-fleet/.env exists (written at enrollment).

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ollie-fleetctl heartbeat
# Headroom for build_health worst case (hung docker/systemctl all timing out)
# so a degraded report still gets posted instead of the unit being killed.
TimeoutStartSec=180
