# templates/systemd/fleet-heartbeat.timer
[Unit]
Description=Run the Ollie Fleet heartbeat every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
RandomizedDelaySec=30
Persistent=false

[Install]
WantedBy=timers.target
